# dist/init/linux-upstart/caddy.conf
description "Caddy HTTP/2 web server"
author "Mathias Beke"

start on runlevel [2345]
stop on runlevel [016]

console logged

setuid www-data
setgid www-data

respawn
respawn limit 10 5

# Letsencrypt-issued certificates will be written to this directory.
env HOME=/etc/ssl/caddy

limit nofile 1048576 1048576

script
	mkdir -p /etc/caddy
	cd /etc/caddy
	exec /usr/bin/caddy -agree=true -email=mail@example.com -log=stdout -conf=/etc/caddy/Caddyfile -root=/var/tmp
end script
